# exchange/supervisor/DeductBalanceQueue.conf
#扣除爆仓余额
[program:queue-deductBalance]
process_name=%(program_name)s
directory=/www/wwwroot/mucw/exchange
command=php artisan queue:work --queue=deductBalance --sleep=1 --tries=3
autostart=true
autorestart=true
user=www
numprocs=1
redirect_stderr=true
stdout_logfile=/www/wwwroot/deductBalance.log
