# exchange/supervisor/leverCloseOrderWorker.conf
#上币脚本队列
[program:leverCloseOrderWorker]
process_name=%(process_num)s
directory=/www/wwwroot/mucw/exchange
command=php artisan lever:closeOrderworker restart
autostart=true
autorestart=true
user=www
numprocs=1
redirect_stderr=true
stdout_logfile=/www/wwwroot/mucw/exchange/storage/logs/closeOrderWorker.log
#workerman程序必须加上这一行,守护停止时也同时停止子进程
stopasgroup=true
