# exchange/supervisor/matchEngine.conf
#撮合引擎0
[program:queue-matchEngine]
process_name=%(program_name)s_%(process_num)02d
directory=/www/wwwroot/test_exchange/exchange
command=php artisan queue:work --queue=matchEngine:%(process_num)1d  --sleep=1 --tries=3
autostart=true
autorestart=true
user=www
numprocs=1
redirect_stderr=true
stdout_logfile=/dev/null
