# exchange/supervisor/socket.conf
#异步改变余额
[program:socket-run]
process_name=%(program_name)s
directory=/www/wwwroot/test_exchange/exchange
command=php artisan socket:run restart
autostart=true
autorestart=true
user=www
numprocs=1
redirect_stderr=true
stdout_logfile=/www/wwwroot/test_exchange/supervisor.log
#workerman程序必须加上这一行,守护停止时也同时停止子进程
stopasgroup=true
